# vrf/components/re4ctor-ipc/deploy/r4-server.socket
[Unit]
Description=R4 Secure RNG socket

[Socket]
ListenStream=/run/r4sock/r4.sock
SocketMode=0660
SocketUser=root
SocketGroup=r4users
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
